# deploy/labmon.service
[Unit]
Description=LabMon GPU monitor
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/opt/labmon
EnvironmentFile=/etc/labmon/labmon.env
ExecStart=/opt/labmon/.venv/bin/uvicorn labmon.app:app --host ${LABMON_HOST} --port ${LABMON_PORT}
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=labmon

[Install]
WantedBy=multi-user.target
